# files/agar.service
[Unit]
Description=AgAR Service

[Service]
ExecStart=/home/pi/agar.sh

[Install]
WantedBy=multi-user.target  
